# vendor/pkgr/scaling/systemd/zammad-web-PROCESS_NUM.service
[Unit]
StopWhenUnneeded=true
Requires=zammad-web.service
After=zammad-web.service

[Service]
Environment=APP_PROCESS_INDEX=__PROCESS_INDEX__
ExecStart=/usr/bin/zammad run web
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n

[Install]
WantedBy=zammad-web.service
