# vendor/pkgr/scaling/systemd/zammad-console-PROCESS_NUM.service
[Unit]
StopWhenUnneeded=true
Requires=zammad-console.service
After=zammad-console.service

[Service]
Environment=APP_PROCESS_INDEX=__PROCESS_INDEX__
ExecStart=/usr/bin/zammad run console
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n

[Install]
WantedBy=zammad-console.service
